# autostart/sparky-screen.desktop
[Desktop Entry]
Type=Application
Name=SparkyScreen
Exec=/usr/bin/python3 /home/pi/sparky-screen/main.py